# root/etc/services.d/mountsize/run
# #!/usr/bin/with-contenv bash
# # shellcheck shell=bash
# function log() {
    # echo "[Mount] ${1}"
# }
# # Copyright (c) 2020, MrDoob
# # All rights reserved.
# # shellcheck disable=SC2086
# umask 022
# mountsize="/app/mountsize/mountsize.sh"
# sleep 5
# log "-> Starting mount size printing <-"
# exec s6-setuidgid abc $mountsize
